# resources/systemd/oc_bifrost.service
[Unit]
Description=Redborder oc_bifrost

[Service]
ExecStart=/usr/lib/redborder/bin/rb_oc_bifrost_start.sh
Restart=always

[Install]
WantedBy=multi-user.target

