# resources/systemd/opscode-erchef.service
[Unit]
Description=Redborder opscode-erchef

[Service]
ExecStart=/usr/lib/redborder/bin/rb_opscode-erchef_start.sh
Restart=always
User=opscode
Group=opscode

[Install]
WantedBy=multi-user.target

